# dund.upstart
description "Bluetooth Dial-Up-Networking Daemon. Provides PPP over RFCOMM services."

start on pld.sysinit-done
stop on pld.shutdown-started

script
	[ -f /etc/sysconfig/bluetooth ] && . /etc/sysconfig/bluetooth

	/usr/bin/dund ${DUND_OPTIONS}
end script
