# pand.upstart
description "Bluetooth Personal Area Networking Daemon. Provides network services over Bluetooth."

start on pld.sysinit-done
stop on pld.shutdown-started

script
	[ -f /etc/sysconfig/bluetooth ] && . /etc/sysconfig/bluetooth

	/usr/bin/pand ${PAND_OPTIONS}
end script
