# rfcomm.upstart
description "Bluetooth RFCOMM setup. Sets up serial devices over Bluetooth."

start on pld.sysinit-done
stop on pld.shutdown-started

script
	/usr/bin/rfcomm -f /etc/bluetooth/rfcomm.conf bind all
end script

post-stop script
	/usr/bin/rfcomm release all
end script
